# etc/cron.hourly/check_bgp
#!/bin/bash
echo "show protocols" | birdc | grep BGP | grep -v Established | mail -Es BGP-report root
