# etc/cron.hourly/watchdog
#!/bin/bash
if ! ps ax | grep -v grep | grep keepalived >/dev/null; then echo "keepalived NOT running - please fix" | mail -s "WATCHDOG: $HOSTNAME" root; fi
